# scripts/bin/systemd/user/ssh-relay.service
[Unit]
Description=SSH Agent Relay for Windows
Documentation=man:ssh-agent(1)
After=network.target

[Service]
Type=simple
ExecStart=%h/.local/bin/ssh-relay.sh
Restart=on-failure
RestartSec=5s

Nice=10
CPUQuota=10%

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
